# misc/42home-fix.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>bricewge.42home-fix</string>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>/Users/bwaegene/.brew/bin:/bin:/usr/bin:/usr/local/bin</string>
    </dict>

    <key>ProgramArguments</key>
    <array>
      <string>/Users/bwaegene/.local/bin/42home-fix.py</string>
    </array>

    <key>OnDemand</key>
    <false/>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/tmp/42home-fix.err</string>
    <key>StandardOutPath</key>
    <string>/tmp/42home-fix.out</string>

  </dict>
</plist>
