# bspwm/.config/systemd/user/polybar@.service
[Unit]
Description=Poybar, create a status bar per monitor
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/run/current-system/sw/bin/polybar --reload example
Environment=MONITOR=%i
Environment=PATH=%h/.local/bin:%h.nix-profile/bin/:/run/current-system/sw/bin

[Install]
WantedBy=graphical-session.target
